# openplotterNetwork/Network/udev/rules.d/11-openplotter-usb0.rules
KERNEL=="usb0", SUBSYSTEMS=="net", RUN+="/bin/bash ~/.openplotter/Network/11-openplotter-usb0.sh"
